# wg-backup.service
[Unit]
Description=WireGuard Telegram Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/wg-backup/wg_backup.py
